# poc/orig_codebase_obsidian_web_mcp/scripts/launchd/com.example.cloudflared-vault.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.cloudflared-vault</string>
    <key>ProgramArguments</key>
    <array>
        <string>REPLACE_WITH_CLOUDFLARED_PATH</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>REPLACE_WITH_HOME/.cloudflared/config-vault-mcp.yml</string>
        <string>run</string>
        <string>vault-mcp</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>REPLACE_WITH_HOME/Library/Logs/cloudflared-vault.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_HOME/Library/Logs/cloudflared-vault-error.log</string>
</dict>
</plist>
